# systemd/redriva.service
# Unité systemd pour Redriva (mode utilisateur)
# À installer manuellement avec :
# systemctl --user enable redriva.service
# systemctl --user start redriva.service

[Unit]
Description=Service Redriva API (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/${USER}/scripts/redriva
ExecStart=/home/${USER}/scripts/redriva/scripts/start.sh
ExecStop=/home/${USER}/scripts/redriva/scripts/stop.sh
Restart=on-failure
User=%i
Environment=VIRTUAL_ENV=/home/${USER}/scripts/redriva/venv
Environment=PYTHONPATH=src

[Install]
WantedBy=default.target
